# etc/services.d/php-error/run
#!/bin/execlineb -P
socat
open:/var/run/php-error.sock
stdout
