# etc/rc.local
#!/bin/sh

# Descendant images may overwrite this to run initialisation tasks at startup
