# etc/services.d/php-fpm/run
#!/bin/sh
echo "init[services]: Starting 'php-fpm' as user 'root'"
php-fpm
