# etc/services.d/nginx/finish
#!/bin/execlineb -S0
# stop the container if this service dies
s6-svscanctl -t /etc/services.d
