# etc/services.d/nginx/run
#!/bin/sh
echo "init[services]: Starting 'nginx' as user 'root'"
nginx
