# etc/services.d/.s6-svscan/finish
#!/bin/execlineb -P
echo "init[shutdown]: Terminating container"
